# timers/backupper-daily.timer
[Unit]
Description=Run a daily incremental backup.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
